# .crontab
@hourly ~/bin/brewup.sh > /dev/null 2>&1
